# deploy/sparkyctrl.service
[Unit]
Description=sparkyctrl worker
After=network-online.target
Wants=network-online.target

[Service]
# Adjust --fence/--token/--audit per host. Bind to the LAN address.
ExecStart=/usr/local/bin/sparkyctrl serve --addr 0.0.0.0:7766 --audit /var/log/sparkyctrl-audit.log
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
